# detect_mobile/detect.service
[Unit]
Description=SecVision
After=syslog.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -u /usr/local/bin/secvision/detect.py --threshold 0.7
Restart=always
RestartSec=120
User=jetsonman
SyslogIdentifier=secvision
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
